# whisper-service/whisper-service.service
[Unit]
Description=Whisper Transcription Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/mindshifting/whisper-service

# Environment variables
Environment="WHISPER_MODEL=base"
Environment="WHISPER_DEVICE=cpu"
Environment="WHISPER_COMPUTE_TYPE=int8"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="CACHE_ENABLED=true"
Environment="CACHE_TTL=3600"
Environment="API_HOST=0.0.0.0"
Environment="API_PORT=8000"
Environment="API_WORKERS=2"
EnvironmentFile=-/opt/mindshifting/whisper-service/.env

# Command
ExecStart=/opt/mindshifting/whisper-service/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
